# crontab
SHELL=/bin/bash
*/20 * * * * /usr/local/bin/python3.10 /andross-api/update.py -d >> /var/log/cron.log 2>&1
59 0 * * * /usr/local/bin/python3.10 /andross-api/update.py -l >> /var/log/cron.log 2>&1